# base-tiny/stuff/etc/inittab
# /etc/inittab: init configuration for SliTaz GNU/Linux.
# Boot-time system configuration/initialization script.
#
::sysinit:/etc/init.d/rcS

# getty respawn shell invocations for selected ttys.
tty1::askfirst:getty 38400 tty1
tty2::askfirst:getty 38400 tty2

# Stuff to do when restarting the init
# process, or before rebooting.
::restart:/etc/init.d/rc.shutdown
::restart:/init
::ctrlaltdel:reboot
::shutdown:/etc/init.d/rc.shutdown

